# playbooks/roles/rvprx/files/main_certbot.timer
[Unit]
Description=Run certbot every day - Generated by Ansible
Documentation=https://github.com/agrimal/minidatacenter

[Timer]
OnCalendar=*-*-* 23:00:00
RandomizedDelaySec=3500
Persistent=true

[Install]
WantedBy=timers.target
