# playbooks/roles/rvprx/files/main_certbot.service
[Unit]
Description=Certbot - Generated by Ansible
Documentation=file:///usr/share/doc/python-certbot-doc/html/index.html
Documentation=https://letsencrypt.readthedocs.io/en/latest/
Documentation=https://github.com/agrimal/minidatacenter

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot -q renew --renew-hook "systemctl restart nginx.service"
PrivateTmp=true
