# playbooks/roles/rvprx/files/restart_service@.service
[Unit]
Description=Restart service %i when file has changed - Generated by Ansible
Documentation=https://github.com/agrimal/minidatacenter

[Service]
Type=oneshot
ExecStart=/usr/bin/env systemctl restart %i.service
