# playbooks/roles/rvprx/templates/restart_service@.path.j2
[Unit]
Description=Restart service %i when file /etc/ssl/private/{{ item.certificate }}-privkey.pem has changed - Generated by Ansible
Documentation=https://github.com/agrimal/minidatacenter

[Path]
PathChanged=/etc/ssl/private/{{ item.certificate }}-privkey.pem
