# playbooks/roles/rvprx/templates/restart_service@.path.j2
[Unit]
Description=Restart service %i when file /etc/ssl/private/{{ (letsencrypt.domains | selectattr('cert_id', 'equalto', item.cert_id) | list).0.domains.0 }}-privkey.pem has changed - Generated by Ansible
Documentation=https://github.com/agrimal/minidatacenter

[Path]
PathChanged=/etc/ssl/private/{{ (letsencrypt.domains | selectattr('cert_id', 'equalto', item.cert_id) | list).0.domains.0 }}-privkey.pem
