# playbooks/roles/rvprx/templates/sync-certs-from-rvprx@.service.j2
[Unit]
Description=Synchronize letsencrypt certificates from %i to this host - Generated by Ansible
Documentation=https://github.com/agrimal/minidatacenter

[Service]
Type=oneshot
ExecStart=/usr/bin/rsync -a --delete-after root@%i:/etc/letsencrypt {{ minidc_working_directory }}/%i/etc/
