# playbooks/roles/rvprx/templates/sync-certs-CERTIFICATE-to-container@.service.j2
[Unit]
Description=Synchronize letsencrypt certificates from this host to the container %i - Generated by Ansible
Documentation=https://github.com/agrimal/minidatacenter

[Service]
Type=oneshot
ExecStart=/usr/bin/rsync --copy-links --backup --suffix "-$(date '+%%F-%%T')" --inplace --times --owner --group --chown={{ item.chown }} --chmod=0644 {{ minidc_working_directory }}/{{ inventory_hostname }}/etc/letsencrypt/live/{{ (letsencrypt.domains | selectattr('cert_id', 'equalto', item.cert_id) | list).0.domains.0 }}/fullchain.pem root@%i:{{ ca_cert_dir }}/{{ (letsencrypt.domains | selectattr('cert_id', 'equalto', item.cert_id) | list).0.domains.0 }}-fullchain.pem
ExecStart=/usr/bin/rsync --copy-links --backup --suffix "-$(date '+%%F-%%T')" --inplace --times --owner --group --chown={{ item.chown }} --chmod=0600 {{ minidc_working_directory }}/{{ inventory_hostname }}/etc/letsencrypt/live/{{ (letsencrypt.domains | selectattr('cert_id', 'equalto', item.cert_id) | list).0.domains.0 }}/privkey.pem root@%i:{{ ca_private_dir }}/{{ (letsencrypt.domains | selectattr('cert_id', 'equalto', item.cert_id) | list).0.domains.0 }}-privkey.pem
