# playbooks/roles/rvprx/templates/sync-certs-CERTIFICATE-to-container@.service.j2
[Unit]
Description=Synchronize letsencrypt certificates from this host to the container %i - Generated by Ansible
Documentation=https://github.com/agrimal/minidatacenter

[Service]
Type=oneshot
ExecStart=/usr/bin/rsync --copy-links --backup --suffix "-$(date '+%%F-%%T')" --inplace --times --owner --group --chown=root:ssl-cert --chmod=0644 {{ minidc_working_directory }}/{{ inventory_hostname }}/etc/letsencrypt/live/{{ item.certificate }}/fullchain.pem root@%i:/etc/ssl/certs/{{ item.certificate }}-fullchain.pem
ExecStart=/usr/bin/rsync --copy-links --backup --suffix "-$(date '+%%F-%%T')" --inplace --times --owner --group --chown=root:ssl-cert --chmod=0640 {{ minidc_working_directory }}/{{ inventory_hostname }}/etc/letsencrypt/live/{{ item.certificate }}/privkey.pem root@%i:/etc/ssl/private/{{ item.certificate }}-privkey.pem
